# reminder_email.crontab
0 13 * * * python /path/to/manage.py send_reminder_emails 2